# maps/reference_app_config_tests.test_supervisord/example.conf
[program:example]
command = /usr/lib/yandex/maps/pycare/enabled/example
user = www-data
umask = 022
environment = PYCR_BIND="unix:/var/run/yandex/maps/pycare/example.sock",PYCR_STATS_SOCK_DIR="/var/run/yandex/maps/pycare/stats"
startretries = 1
stopsignal = TERM
stopwaitsecs = 120
stdout_logfile_maxbytes = 1MB
stderr_logfile_maxbytes = 1MB
stdout_logfile_backups = 2
stderr_logfile_backups = 2
autorestart = true
autostart = false
