# maps/reference_app_config_tests.test_generate_supervisord_/yacare-testapp.conf
[program:yacare-testapp]
command = /usr/bin/env YCR_MODE="fastcgi:/var/run/yandex/maps/yacare/yacare-testapp.sock" /usr/lib/yandex/maps/yacare/enabled/yacare-testapp
user = www-data
umask = 022
startsecs = 0
stopsignal = TERM
stopwaitsecs = 120
stdout_logfile_maxbytes = 1MB
stderr_logfile_maxbytes = 1MB
stdout_logfile_backups = 2
stderr_logfile_backups = 2
autostart = false
autorestart = true
