# classfields/sphinx-indexer-test/files/etc/cron.d/autoru-sphinx-indexer
# Index Boards
*/10 * * * * nobody zk-flock autoru_boards-reindex "/usr/bin/indexer -c  /etc/sphinx/conf.d/boards.conf *_delta --rotate" > /var/log/sphinx/indexer_boards-reindex-delta1.log 2>&1

# Index boards all
#0 2 1 * * nobody zk-flock autoru_boards-reindex-all "/usr/bin/indexer -c  /etc/sphinx/conf.d/boards.conf --all --rotate" >> /var/log/sphinx/indexer_boards-reindex-all.log 2>&1

# office7
*/2 7-20 * * * nobody zk-flock autoru_office7-reindex_parts "/usr/bin/indexer -c  /etc/sphinx/conf.d/office7.conf office7_clients2 office7_agents2 --rotate" >> /var/log/sphinx/indexer_office7-reindex-parts.log 2>&1
