# classfields/sphinx-indexer-test/files/etc/cron.d/autoru-sphinx-moderation
0 21 * * 2 sphinx /usr/bin/indexer --config /etc/sphinx/conf.d/moderation.conf all7_plain_1 --rotate > /var/log/sphinx/indexer_moderation_all7_1.log  2>&1
0 21 * * 2 sphinx /usr/bin/indexer --config /etc/sphinx/conf.d/moderation.conf all7_plain_2 --rotate > /var/log/sphinx/indexer_moderation_all7_2.log  2>&1
0 21 * * 2 sphinx /usr/bin/indexer --config /etc/sphinx/conf.d/moderation.conf all7_plain_3 --rotate > /var/log/sphinx/indexer_moderation_all7_3.log  2>&1
0 21 * * 2 sphinx /usr/bin/indexer --config /etc/sphinx/conf.d/moderation.conf all7_plain_4 all7_plain_5 --rotate > /var/log/sphinx/indexer_moderation_all7_4.log 2>&1
10 15 * * 7 sphinx /usr/bin/indexer --config /etc/sphinx/conf.d/moderation.conf users8_plain --rotate > /var/log/sphinx/indexer_moderation_users_1.log  2>&1

