# classfields/cassandra-backup-s3-testing/files/etc/cron.d/cassandra-backup-s3-testing
PATH=/usr/bin:/usr/sbin:/sbin:/bin

0 16 * * * root bash /usr/sbin/cassandra-backup-s3.sh 2>&1 | /usr/bin/logger -t cassandra-backup-s3
