# billing/balance-test-apikeys-uwsgi.conf
description "Apikeys TEST backend via uwsgi"
start on runlevel [2345] and local-filesystems and net-device-up IFACE=lo
stop on runlevel [016]

respawn limit unlimited

env YANDEX_XML_CONFIG=/etc/yandex/balance-test-apikeys/servant.nonprod.cfg.xml
exec /opt/yandex-balance-apikeys/balance_test_apikeys_uwsgi_app /etc/yandex/balance-test-apikeys/balance-test-apikeys-uwsgi.ini --die-on-term

pre-start script
    mkdir -p -m 0775 /var/run/yb
    chown yb:www-data /var/run/yb
end script
