# billing/test_xmlrpc (2).conf
[program:test_xmlrpc]
priority=20
command=/binaries/balance-binary test-xmlrpc
environment=YANDEX_XML_CONFIG="/tools/test_xmlrpc/servant-docker.cfg.xml"
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
