# billing/testserv (2).conf
[program:testserv]
priority=20
command=/opt/venvs/apikeys/bin/uwsgi --processes 4 --enable-threads --master --http-socket :8025 --plugin python --wsgi-file test_apikeys_servant/balance_test_apikeys_uwsgi_app.py
directory=/apikeys
environment=YANDEX_XML_CONFIG="/apikeys/test_apikeys_servant/servant-debug.cfg.xml",PYTHONPATH="/apikeys"
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
