# billing/testserv.conf
[program:testserv]
priority=20
command=/opt/yandex-balance-apikeys/balance_test_apikeys_uwsgi_app --processes 4 --enable-threads --master --http-socket :8025 -w billing.apikeys.test_apikeys_servant.balance_test_apikeys_uwsgi_app
environment=YANDEX_XML_CONFIG="/configs/balance-test-apikeys-uwsgi.dbg.cfg.xml"
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
