# billing/test_xmlrpc.conf
[program:test_xmlrpc]
priority=20
command=python test_xmlrpc_servant.py
directory=/tools/test_xmlrpc
environment=YANDEX_XML_CONFIG="/tools/test_xmlrpc/servant-docker.cfg.xml",PYTHONPATH="/tools"
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
