# billing/test_xmlrpc/upstart/balance-test-xmlrpc.conf
description "Primary Balance XML-RPC API for Yandex services"

start on runlevel 2
stop on runlevel [!2]

respawn
respawn limit unlimited

env SERVANT_NAME="balance-test-xmlrpc"

# Only 15 symbols are allowed
env PROCESS_NAME="yb-test-xmlrpc"

script
    set +e

    . /etc/oraprofile

    ENV_CONFIG="/etc/yandex/balance-common/environment.cfg.xml"
    ENV_TYPE=`xpath -q -e '/*/Type/text()' "$ENV_CONFIG"`
    CONFIG_DIR="/etc/yandex/$SERVANT_NAME"
    if [ "$ENV_TYPE" = prod ] ; then
        echo "TestXmlRpc is not for use in production"
        exit 5
    fi
    export YANDEX_XML_CONFIG="$CONFIG_DIR/servant.nonprod.cfg.xml"
    if [ "$ENV_TYPE" = prod ] || [ ! -f "$YANDEX_XML_CONFIG" ] ; then
        YANDEX_XML_CONFIG="$CONFIG_DIR/servant.cfg.xml"
    fi
    export OMNIORB_CONFIG="$CONFIG_DIR/omniORB.nonprod.cfg"
    if [ "$ENV_TYPE" = prod ] || [ ! -f "$OMNIORB_CONFIG" ] ; then
        OMNIORB_CONFIG="$CONFIG_DIR/omniORB.cfg"
    fi

    USER="yb"
    LOG_FILE=`xpath -q -e '/*/Log[contains(., ".log")][1]/text()' "$YANDEX_XML_CONFIG"`
    PID_FILE=`xpath -q -e '/*/Pid/text()' "$YANDEX_XML_CONFIG"`
    SERVANT_BIN="/usr/bin/$PROCESS_NAME"

    if [ ! -d "/var/log/$USER" ] ; then
        mkdir -p "/var/log/$USER"
        chown "$USER:$USER" "/var/log/$USER"
    fi
    if [ ! -d "/var/log/$USER/arc" ] ; then
        mkdir -p "/var/log/$USER/arc"
        chown "$USER:$USER" "/var/log/$USER/arc"
    fi
    if [ ! -d "/var/run/$USER" ] ; then
        mkdir -p "/var/run/$USER"
        chown "$USER:$USER" "/var/run/$USER"
    fi
    touch "$LOG_FILE"
    chown "$USER:$USER" "$LOG_FILE"
    exec start-stop-daemon \
        --start \
        --chuid "$USER" \
        --pidfile "$PID_FILE" \
        --name "$PROCESS_NAME" \
        --exec "$SERVANT_BIN" >> "$LOG_FILE" 2>> "$LOG_FILE" -- "test-xmlrpc"
end script

post-start script
    set +e

    ENV_CONFIG="/etc/yandex/balance-common/environment.cfg.xml"
    ENV_TYPE=`xpath -q -e '/*/Type/text()' "$ENV_CONFIG"`
    CONFIG_DIR="/etc/yandex/$SERVANT_NAME"
    export YANDEX_XML_CONFIG="$CONFIG_DIR/servant.nonprod.cfg.xml"
    if [ "$ENV_TYPE" = prod ] || [ ! -f "$YANDEX_XML_CONFIG" ] ; then
        YANDEX_XML_CONFIG="$CONFIG_DIR/servant.cfg.xml"
    fi

    HOST="127.0.0.1"
    PORT=`xpath -q -e '/*/TestXmlRpc/HttpListener/Port/text()' "$YANDEX_XML_CONFIG"`

    for i in `seq 1 40`; do
        /usr/bin/yb-servant-ping "$HOST:$PORT"

        if [ $? = 0 ]; then
            exit 0
        fi
        sleep 1
    done

    stop
end script

