# pay/it_tests/dockered_hudsucker/deploy/etc/supervisor/conf.d/graphite.conf
[program:graphite]
command=python3 /opt/graphite.py
stdout_logfile=/var/log/supervisor/graphite.log
autostart=true
autorestart=true
